# systemd/dsl-buttons.service
[Unit]
Description=DSL door buttons
After=redis.service
StopWhenUnneeded=true

[Service]
Restart=always
User=root
ExecStart=/usr/bin/dsl-buttons
