# systemd/hackman-rfidd.service
[Unit]
Description=Run hackman rfid daemon
Wants=redis-server.service

[Service]
WorkingDirectory=/var/www/hackman/
EnvironmentFile=/etc/hackman.env
Restart=always
User=root
ExecStart=python3 manage.py rfidd

[Install]
WantedBy=multi-user.target
