# systemd/hackman-rfidd.service
[Unit]
Description=Run hackman rfid daemon
Wants=redis-server.service

[Service]
EnvironmentFile=/etc/hackman.env
Restart=always
User=root
ExecStart=/usr/bin/hackman-manage rfidd

[Install]
WantedBy=multi-user.target
