# systemd/hackman-doord.service
[Unit]
Description=Run hackman door daemon
Wants=redis-server.service

[Service]
WorkingDirectory=/var/www/hackman/
EnvironmentFile=/etc/hackman.env
Restart=always
User=root
ExecStart=python3 manage.py doord

[Install]
WantedBy=multi-user.target
